# tools/log_db/systemd/vps-logdb.service
[Unit]
Description=VPS Log DB (PostgreSQL 16, port 5433)
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c '/usr/bin/pg_ctlcluster 16 main start || /usr/bin/pg_ctlcluster 16 main status'
ExecStop=/usr/bin/pg_ctlcluster 16 main stop

[Install]
WantedBy=default.target
